# ops/backup/crumb-db-backup.service
[Unit]
Description=CrumbVMS nightly PostgreSQL backup to the archive disk
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/crumb/scripts/db-backup.sh
