# ops/backup/crumb-db-backup.timer
[Unit]
Description=Nightly CrumbVMS DB backup (03:30 America/Los_Angeles)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
